# eip/LogstashInstaller/linux/scripts/cynet_logstash.service
[Unit]
Description=Cynet Logstash Service

[Service]
ExecStart=#CMD
Restart=always

[Install]
RequiredBy=network.target